# jobagent.service
[Unit]
Description=Job Agent v2.0 - LinkedIn Job Notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/job-agent
EnvironmentFile=/home/ubuntu/job-agent/.env
ExecStart=/home/ubuntu/job-agent/venv/bin/python3 agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=job-agent

[Install]
WantedBy=multi-user.target
